# front-end/supervisord.conf
[supervisord]
nodaemon=true

[program:frontend]
command=npm run dev
autostart=true
autorestart=true
startretries=3
                